# rebeccablacklinux_files/usr/import/etc/systemd/system/multi-user.target.wants/rngd.service
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
ExecStart=/sbin/rngd -f

[Install]
WantedBy=multi-user.target
