# rebeccablacklinux_files/usr/import/lib/systemd/system/packagekit.service
[Unit]
Description=PackageKit Daemon

[Service]
Type=dbus
BusName=org.freedesktop.PackageKit
User=root
ExecStart=/usr/lib/packagekit/packagekitd
