# rebeccablacklinux_files/usr/import/lib/systemd/system/usbmuxd.service
[Unit]
Description=iOS USB Muxer

[Service]
Type=simple
ExecStart=/usr/sbin/usbmuxd -u -U usbmuxd
PIDFile=/var/run/usbmuxd.pid
